# local-modules/eizo/files/firewall/firewall.service
[Unit]
Description = Setup networking and firewall

[Service]
Type = oneshot
ExecStart = /etc/firewall/run
