# local-modules/eizo/files/ddns/ddns-updater.service
[Unit]
Description = ddns-updater
After = network.target

[Service]
Type = oneshot
ExecStart = /usr/local/bin/ddns-updater
