# local-modules/eizo/files/ddns/ddns-updater.timer
[Unit]
Description=Run ddns-updater at regular interval

[Timer]
OnBootSec=1min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
