# local-modules/eizo/files/transmission/transmission-daemon.service
[Service]
ExecStart = /usr/bin/transmission-daemon -f --log-error --config-dir=/var/lib/transmission-daemon/info
Restart = on-failure
