# local-modules/eizo/files/kodi/kodi.service
[Unit]
Description = Kodi Media Center
After = remote-fs.target
After = syslog.target
 
[Service]
Type = simple
Restart = always
RestartSec = 15
User = xbmc
Group = xbmc
ExecStart = /usr/local/bin/kodi-start
ExecStop = /usr/local/bin/kodi-stop
TTYPath = /dev/tty9
StandardInput = tty
StandardOutput = journal

[Install]
WantedBy = multi-user.target
